# docs/devel.crontab
*/5     *       *       *       *       root    /usr/local/bin/curl -s http://mystery.vmkteam.ru/int/daemons/sync-sources/ > /dev/null
0-59     *       *       *       *       root    /usr/local/bin/curl -s http://mystery.vmkteam.ru/int/daemons/process-queue/ > /dev/null
0-59     *       *       *       *       root    /usr/local/bin/curl -s http://mystery.vmkteam.ru/int/daemons/process-queue/ > /dev/null
0-59     *       *       *       *       root    /usr/local/bin/curl -s http://mystery.vmkteam.ru/int/daemons/process-queue/ > /dev/null
0-59     *       *       *       *       root    /usr/local/bin/curl -s http://mystery.vmkteam.ru/int/daemons/process-queue/ > /dev/null
0-59     *       *       *       *       root    /usr/local/bin/curl -s http://mystery.vmkteam.ru/int/daemons/process-queue/ > /dev/null